# dashboard/k3s-dashboard-ui.service
[Unit]
Description=K3s Dashboard UI Web Server (Python)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/k3s-dashboard
ExecStart=/usr/bin/python3 -m http.server 8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
